# support/pkg/arch/eruption.install
post_install() {
    systemctl preset eruption
    systemctl reload dbus.service
    udevadm control --reload-rules
    systemctl try-restart eruption.service || udevadm trigger -s usb
}

pre_remove() {
    systemctl stop eruption.service
}

post_remove() {
    systemctl daemon-reload
}
